# repository/ssh-sync-daemon/ltpda-ssh-sync.service
[Unit]
Description=LTPDA SSH Sync Daemon
Documentation=https://github.com/your-org/ltpda
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/ltpda-ssh-sync/ssh_sync_daemon.py
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ltpda-ssh-sync

[Install]
WantedBy=multi-user.target
